# supervisor/kafka1.conf
[program:kafka1]
command=/usr/bin/start-kafka-n.sh 1
autostart=true
autorestart=true

