# supervisor/kafka0.conf
[program:kafka0]
command=/usr/bin/start-kafka-n.sh 0
autostart=true
autorestart=true

